# docker/php/crontab
@hourly /usr/local/bin/php /var/www/html/sf4/bin/console gpao:cron >> /var/log/cron.log 2>&1
# This extra line makes it a valid cron